# docker/rootfs/etc/services.d/mqtt_client/run
#!/bin/bash

chmod +x /fevr/run_mqtt_client.sh
/fevr/run_mqtt_client.sh