# docker/rootfs/etc/services.d/mqtt_client/run
#!/bin/bash

#Check if Virtual Environment is initialized
/fevr/run_mqtt_client.sh