# docker/fEVR/rootfs/etc/services.d/mqtt_client/run
#!/bin/bash

# Get Environment Variables
/opt/fevr/startup/mqtt_client